# modules/desktop/hyprland/assets/atlas-hypr-check.service
[Unit]
Description=Atlas — check whether Hyprland is installable yet
Documentation=file:%h/atlas/modules/desktop/hyprland/README.md

[Service]
Type=oneshot
ExecStart=%h/.local/bin/atlas-hypr-check.sh
