# modules/desktop/hyprland/assets/atlas-hypr-check.timer
[Unit]
Description=Atlas — daily Hyprland availability check

[Timer]
OnBootSec=10min
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
